# infra/ec2/auto-shutdown.service
[Unit]
Description=ff-training auto-shutdown idle check
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ff-auto-shutdown
